# etc/systemd/snif-conn.service
[Unit]
Description=SNIF Connector

[Service]
User=snif
Group=snif
ExecStart=/usr/bin/snif-conn -d

[Install]
WantedBy=multi-user.target
