# etc/systemd/snif-certd.service
[Unit]
Description=SNIF CA Proxy

[Service]
User=snif
Group=snif
ExecStart=/home/snif/bin/snif-certd

[Install]
WantedBy=multi-user.target
